# services/rclonedash.service
[Unit]
Description=RcloneDash - Interface Web pour rclone
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/RcloneDash
ExecStart=/usr/bin/python3 %h/.local/share/RcloneDash/rclone-monitor.py
Restart=on-failure

[Install]
WantedBy=default.target
