# services/rclone-bisync.service.template
[Unit]
Description=Synchronisation rclone Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-__HOME__/.config/rclone/bwlimit.env
# La garde ne lance le bisync que si c'est utile (changement local récent,
# sync complet périodique, ou déclenchement manuel). Voir rclone-bisync-guard.sh.
ExecStart=__HOME__/.local/share/RcloneDash/rclone-bisync-guard.sh
